# playbooks/hermes/templates/hermes-dashboard.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=Hermes Agent dashboard (web UI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Serve the prebuilt UI dist (built during the install phase) with --skip-build so
# the unit never reruns tsc+vite at start. Bind non-loopback WITHOUT --insecure so
# the native scrypt BasicAuthProvider gate (HERMES_DASHBOARD_BASIC_AUTH_* in
# ~/.hermes/.env) is enforced — Hermes refuses to bind a public host unless an auth
# provider is registered, which is why this unit is only enabled when auth is set.
ExecStart=%h/.local/bin/hermes dashboard --host {{ hermes_dashboard_host }} --port {{ hermes_dashboard_port }} --skip-build --no-open
Restart=on-failure
RestartSec=5
{% if hermes_dashboard_enable_podman_chat | default(false) | bool %}
# Opt-in: dashboard Chat is spawned by this web UI unit; gateway, Slack, and
# crons use separate service paths, so this only relaxes dashboard hardening.
# The terminal backend image defaults still assume the igou-devenv container
# contract unless the container identity vars and image are changed together.
# Rootless Podman needs writable user runtime state and bind-mount sources.
# It also needs namespace creation and bind-mount preparation/relabeling, so
# RestrictNamespaces=true, NoNewPrivileges=true, and @system-service filtering are
# incompatible with Podman-enabled dashboard Chat. Testing also showed the
# remaining systemd sandboxing directives below, including ReadWritePaths, make
# newuidmap fail with EPERM in this user-service context, so the opt-in branch
# omits them. This weakens dashboard hardening and is intentionally disabled by
# default.
NoNewPrivileges=false
RestrictNamespaces=
SystemCallFilter=
SystemCallErrorNumber=
Environment="HERMES_DOCKER_BINARY=/usr/bin/podman"
{% else %}
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.hermes
PrivateTmp=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
LockPersonality=true
{% endif %}

[Install]
WantedBy=default.target
