# playbooks/grafana-kiosk/templates/cog-kiosk.service.j2
# {{ ansible_managed }}
[Unit]
Description=Grafana kiosk (cog/WPE WebKit on DRM)
Wants=network-online.target
After=network-online.target nginx.service seatd.service
Requires=seatd.service
Conflicts=getty@tty1.service
After=getty@tty1.service

[Service]
User={{ _kiosk_user }}
Group={{ _kiosk_user }}
SupplementaryGroups=video render input
# PAM session on tty1 gives cog seat0 and a runtime dir; seatd hands it DRM.
PAMName=login
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
StandardError=journal
Environment=XDG_SEAT=seat0
Environment=LIBSEAT_BACKEND=seatd
# The PAM login session hands the unit capabilities; WPE's bwrap sandbox
# refuses to start a capable-but-not-setuid process ("Unexpected
# capabilities but not setuid"). The browser needs no caps - drop them all.
CapabilityBoundingSet=
{# Extra env from kiosk_cog_env (e.g. WEBKIT_DISABLE_* on virtual GPUs). #}
{% for k, v in _kiosk_cog_env | dictsort %}
Environment="{{ k }}={{ v }}"
{% endfor %}
{% if kiosk_drm_video_mode | default('') %}
Environment=COG_PLATFORM_DRM_VIDEO_MODE={{ kiosk_drm_video_mode }}
{% endif %}
# --webprocess-failure=restart: 'restart-exit' is not a valid action on the
# cog 0.16 that Bookworm/RPi OS ship (valid: error-page, exit, exit-ok,
# restart) - the unit crash-loops with "Invalid action name". WebProcess
# crashes restart in-place; full cog exits come back via Restart=always.
ExecStart=/usr/bin/cog --platform=drm --webprocess-failure=restart "{{ _kiosk_cog_url }}"
Restart=always
RestartSec=5
# Daily recycle: browsers leak on live dashboards; Restart=always brings the
# session straight back.
RuntimeMaxSec={{ _kiosk_service_max_runtime }}

[Install]
WantedBy=multi-user.target
